# ansible/roles/prometheus/templates/prometheus.service.j2
[Unit]
Description=Prometheus — Sistema de monitoreo y alertas
After=network.target

[Service]
Type=simple
User={{ prometheus_user }}

# ExecStart: comando para iniciar Prometheus. Flags:
#   --config.file       → archivo de configuración principal (targets, alertmanager, reglas)
#   --storage.tsdb.path → directorio de la TSDB (Time Series Database), el motor de
#                         almacenamiento de Prometheus. Optimizado para series de tiempo;
#                         almacena los datos en bloques de 2h que se compactan periódicamente.
#   --storage.tsdb.retention.time → cuánto tiempo conservar datos. 15d ≈ 170MB con 200 series
#                         (estimación: ~2-3 bytes/muestra × 4 muestras/min × 200 series × 15d)
#   --web.listen-address → puerto HTTP donde escucha la UI y la API de Prometheus
#   --web.enable-lifecycle → habilita endpoints de gestión en tiempo real:
#                            POST /-/reload → recarga prometheus.yml sin reiniciar el proceso
#                                            (lo usan los handlers de Ansible)
#                            POST /-/quit   → apaga Prometheus limpiamente
#     ADVERTENCIA: sin autenticación, cualquier cliente con acceso al puerto puede hacer reload.
#     En el taller (SG abierto a 0.0.0.0/0) es un riesgo aceptable.
#     En producción: restringir el SG a IPs del equipo de ops o añadir autenticación básica.
ExecStart={{ install_dir }}/prometheus \
  --config.file={{ prometheus_config_dir }}/prometheus.yml \
  --storage.tsdb.path={{ prometheus_data_dir }} \
  --storage.tsdb.retention.time=15d \
  --web.listen-address=":{{ prometheus_port }}" \
  --web.enable-lifecycle

# ExecReload: comando ejecutado al recibir `systemctl reload prometheus`.
# SIGHUP es la señal Unix tradicional de "recarga configuración".
# Prometheus la interpreta como: vuelve a leer prometheus.yml sin interrumpir scraping.
# $MAINPID es la variable de systemd con el PID del proceso principal.
ExecReload=/bin/kill -HUP $MAINPID

# Restart=on-failure: systemd reiniciará el proceso si termina con código de error > 0.
# Cubre crashes o errores inesperados. No reinicia si se detiene con `systemctl stop`.
Restart=on-failure
# RestartSec: espera 5 segundos antes de reintentar el arranque.
# Evita que un servicio en crash loop consuma todos los recursos del sistema.
RestartSec=5s

[Install]
# WantedBy=multi-user.target: el servicio arranca automáticamente cuando el sistema
# alcanza el modo multiusuario (nivel de ejecución 3 en SysV): modo normal de servidor,
# sin interfaz gráfica. Es el target estándar para servicios de red en Linux.
WantedBy=multi-user.target
