# ansible/roles/alertmanager/templates/alertmanager.service.j2
[Unit]
Description=Alertmanager — Gestión de alertas de Prometheus
After=network.target

[Service]
Type=simple
User={{ alertmanager_user }}

# --storage.path: directorio donde Alertmanager persiste silences (alertas silenciadas
# manualmente desde la UI) y el notification log (historial de qué notificaciones se
# enviaron y cuándo). Sin persistencia, cada reinicio olvidaría los silences activos.
ExecStart={{ install_dir }}/alertmanager \
  --config.file={{ alertmanager_config_dir }}/alertmanager.yml \
  --storage.path={{ alertmanager_data_dir }} \
  --web.listen-address=":{{ alertmanager_port }}"

# Al igual que Prometheus, Alertmanager soporta SIGHUP para recargar configuración.
# Esto permite cambiar las rutas de alertas y receptores sin perder los silences activos.
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
