# ansible/roles/node_exporter/templates/node_exporter.service.j2
[Unit]
Description=Prometheus Node Exporter — métricas del sistema operativo
After=network.target

[Service]
Type=simple
User={{ node_exporter_user }}
ExecStart={{ install_dir }}/node_exporter \
  --web.listen-address=":{{ node_exporter_port }}" \
  --collector.systemd \
  --collector.processes
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
